# scripts/profile.d/xinu.sh
XINU_SERVERS="%SERVERNAME%"
export  XINU_SERVERS
